# cron-docker/alpine/crontab
* * * * * /app/script.sh >> /app/cron.log 2>&1
